# meta-fso/recipes-freesmartphone/freesmartphone/frameworkd/frameworkd.service
[Unit]
Description=FSO1 daemon used for opimd
After=dbus.socket

[Service]
Type=dbus
BusName=org.freesmartphone.frameworkd
ExecStart=/usr/bin/frameworkd

[Install]
WantedBy=basic.target
