# meta-samsung/recipes-core/systemd/systemd-machine-units/tuna/pvrinit.service
[Unit]
Description=Set up device PVR graphics chip

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/vendor/bin/pvrsrvctl --no-module --start

[Install]
WantedBy=basic.target
