# meta-samsung/recipes-core/systemd/systemd-machine-units/mountdevtmpfs.service
[Unit]
Description=Custom mount of devtmpfs for samsung

[Service]
ExecStart=/usr/bin/mountdevtmpfs.sh

[Install]
WantedBy=basic.target
