# meta-samsung/recipes-core/systemd/systemd-machine-units/rndissetup.service
[Unit]
Description=Set up rndis interface on samsung devices

[Service]
ExecStart=/usr/bin/rndissetup.sh

[Install]
WantedBy=basic.target
