# meta-samsung/recipes-core/systemd/systemd-machine-units/rndissetup.service
[Unit]
Description=Set up rndis interface on samsung devices
Before=basic.target
DefaultDependencies=no

[Service]
ExecStart=/usr/bin/rndissetup.sh
