# meta-lg/recipes-core/android-system/android-system/android-system-filesystems.target
[Unit]
Description=Android System Filesystems requirements
Before=android-system.service
After=system.mount persist.mount firmware.mount
Requires=system.mount persist.mount firmware.mount

[Install]
RequiredBy=android-system.service
