# meta-android/recipes-core/android-system/android-system/android_system.service
[Unit]
Description=Android System startup

[Service]
Type=simple
Restart=on-failure
ExecStartPre=/usr/bin/android_system_setup
ExecStart=/sbin/init.android

[Install]
WantedBy=basic.target
