# meta-xiaomi/recipes-core/systemd/systemd-machine-units/tissot/wifi-module-load.service
[Unit]
Description=Load wifi module
After=android-system.service
Conflicts=shutdown.target actdead.target

[Service]
Type=simple
RemainAfterExit=yes
ExecStart=/sbin/modprobe wlan
ExecStop=/sbin/modprobe -r wlan
Restart=on-failure
RestartSec=2

[Install]
WantedBy=basic.target
