# meta-openmoko/recipes-core/udev/udev-extraconf/gps.rules
KERNEL=="ttySAC1", MODE="0660", RUN+="/etc/udev/scripts/gps.sh"
