# meta-samsung/recipes-core/initscripts/android-initscripts/pvrinit.conf
start on stopped finish

script
	/vendor/bin/pvrsrvctl --no-module --start
end script
